# mail/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnActiveSec=20
OnUnitActiveSec=5m
Unit=mbsync.service

[Install]
WantedBy=timers.target
